# services/bee-oled.service
[Unit]
Description=Bee OLED UI (NanoHat OLED)
After=multi-user.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=root
WorkingDirectory=/bee/services
ExecStart=/bee/services/bee_oled.py
Restart=always
RestartSec=2
KillMode=process

[Install]
WantedBy=multi-user.target
