# services/bee-shell.service
[Unit]
Description=Bee Web Shell (uvicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/mnt/data/bee/web
EnvironmentFile=/etc/bee-shell.env
ExecStart=/mnt/data/bee/bee_env/bin/uvicorn server:app --host 0.0.0.0 --port 8088
Restart=always
RestartSec=2
NoNewPrivileges=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target
